# scheduled-job.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
launchd plist template — scheduled-operations job (Job A)
macOS only. Linux uses systemd; this template does not apply there.

To instantiate for a new vault:
  1. Replace {{VAULT_PREFIX}} with the vault's short prefix
     (e.g. sc, is, sps, fin, qf, mktan).
  2. Replace {{VAULT_PATH}} with the absolute path to the vault root
     (e.g. /Users/<username>/vaults/<vault-slug>). launchd does not
     expand ~ or $HOME — the path must be literal.
  2b. Replace {{CLAUDE_BIN}} with the absolute path to the claude
     binary — the output of `which claude`. It varies by install:
     /opt/homebrew/bin/claude (Apple Silicon Homebrew),
     /usr/local/bin/claude (Intel Homebrew), or an npm global bin.
     launchd jobs get a minimal PATH, so the path must be absolute.
  3. Optionally adjust StartCalendarInterval > Hour if you want this
     vault's job to fire at a different time of day. Default 08:00.
  4. Replace {{INTAKE_DIR}} with the absolute path to the intake
     staging directory (e.g. /Users/<username>/Downloads). The
     `--add-dir` argument grants the headless session read/write access
     to it; launchd does not expand ~ or $HOME, so the path must be
     literal. Without it, Intake cannot enumerate
     `<vault-prefix>-*.md` candidates.
     IMPORTANT: `--add-dir` is variadic (it accepts multiple paths), so
     it MUST come AFTER the quoted prompt, not before it. If placed
     before the prompt, the parser swallows the prompt as a second
     directory and claude exits 1 with "Input must be provided ...
     when using --print".
  5. The `--permission-mode bypassPermissions` argument grants the
     headless session non-interactive write approval. Without it,
     headless `claude -p` blocks on Write/Edit/Bash operations
     ("awaiting permission") and the digest and log.md anchor are
     never written. Single-valued flag; position is tolerant but
     conventionally placed between the quoted prompt and `--add-dir`.
  6. Save as
       ~/Library/LaunchAgents/com.user.vault-{{VAULT_PREFIX}}-scheduled.plist
  7. plutil -lint <path>   — verify well-formed XML.
  8. launchctl load <path> — register with launchd.
  9. launchctl list | grep com.user.vault — confirm loaded.

IMPORTANT: the shell command in ProgramArguments uses &amp;&amp; instead
of &&. This is XML entity-encoding for the ampersand character. macOS
launchd's CFPropertyList parser is lenient and accepts unescaped &, but
plutil -lint (strict XML) rejects it. Always use the escaped form.
launchd unescapes the entities back to & before handing the command to
zsh, so the actual shell command is identical either way.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.vault-{{VAULT_PREFIX}}-scheduled</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd {{VAULT_PATH}} &amp;&amp; {{CLAUDE_BIN}} -p "Run scheduled-operations workflow as defined in CLAUDE.md. Apply the 5-day guard." --permission-mode bypassPermissions --add-dir {{INTAKE_DIR}}</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{VAULT_PATH}}/.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>{{VAULT_PATH}}/.launchd.err</string>
</dict>
</plist>
